# templates/default/mirror.crontab.erb
MAILTO=<%= node["mirror"]["mailto"] %>
0 2 * * * /usr/local/bin/mirror > /dev/null
